# debian/sidecar/magma-sidecar.service
[Unit]
Description=Magma Sidecar - tip-based MEV reprioritization for Monad
Documentation=https://github.com/MagmaStaking/magma-sidecar
After=network.target
Wants=network.target

[Service]
Type=simple
User=magma-sidecar
Group=magma-sidecar
# All knobs live here. Every flag in src/config.rs has a 1:1 MAGMA_* env var,
# so ExecStart stays argument-less and operators edit only this file.
EnvironmentFile=/etc/magma-sidecar/sidecar.env
ExecStart=/usr/bin/magma-sidecar
Restart=always
RestartSec=2

# Security hardening. The process has no writable filesystem paths; it connects
# to the ACL-protected txpool socket under /run/monad-ipc and serves
# observability on loopback only.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
CapabilityBoundingSet=
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
ReadOnlyPaths=/

# Only local Unix IPC and loopback observability are permitted. IPAddressDeny
# applies to ingress and egress; explicit loopback allows keep /health and
# /metrics available to local monitoring without permitting network exfiltration.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
IPAddressDeny=any
IPAddressAllow=127.0.0.0/8
IPAddressAllow=::1/128

# Resource containment (cgroup v2) so a compromised or runaway sidecar can
# never starve the validator's CPU/memory/IO. Memory/task caps are conservative
# starting points — retune against measured peak usage (`systemctl status
# magma-sidecar` reports memory peak; `systemd-cgtop`).
MemoryHigh=512M
MemoryMax=512M
MemorySwapMax=0
# Matches the validator-sidecar security profile. Validate thread pressure on
# representative validator hardware before lowering further.
TasksMax=64
# CPU pinning: run on the RPC cores (12-15), not the node's consensus cores
# (8-11). The RPC path is far less latency-critical than consensus, so sharing
# it keeps the reprioritizer off the hot cores at negligible cost. Operators
# whose core layout differs can override these via a drop-in
# (`systemctl edit magma-sidecar`; survives upgrades). See docs/RELEASING.md.
AllowedCPUs=12-15
CPUAffinity=12,13,14,15
# Yield to co-located services under contention (CPUWeight, default 100) and cap
# absolute burn at one core-equivalent (CPUQuota). Light, latency-sensitive.
CPUWeight=50
CPUQuota=100%
# Deprioritize vs the node under disk contention.
IOWeight=50
# FD ceiling: observability HTTP + txpool IPC.
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
