# debian/sidecar/postinst
#!/bin/sh
# postinst — dpkg(1) maintainer script. See `man 7 deb-postinst` for the full
# argument contract; we only handle `configure` (fresh install + upgrade).

set -e

UNIT=magma-sidecar.service
ENV_FILE=/etc/magma-sidecar/sidecar.env
ENV_EXAMPLE=/etc/magma-sidecar/sidecar.env.example

case "$1" in
    configure)
        # Seed the operator env file from the example on **first** install only.
        # `conffiles` already protects sidecar.env.example from being clobbered by
        # operator edits to sidecar.env on subsequent upgrades.
        if [ ! -e "${ENV_FILE}" ] && [ -e "${ENV_EXAMPLE}" ]; then
            cp "${ENV_EXAMPLE}" "${ENV_FILE}"
            chmod 0640 "${ENV_FILE}"
            # If the validator host already has the `monad` group (the same uid
            # that owns the txpool socket), grant it read on the env file so the
            # service can pick up MAGMA_* vars. Silent fallback on dev hosts.
            chgrp monad "${ENV_FILE}" 2>/dev/null || true
        fi

        # Register/reload the unit with systemd, but *do not* enable it
        # automatically — operators must edit /etc/magma-sidecar/sidecar.env
        # (at minimum MAGMA_MONAD_RPC_URL, MAGMA_TXPOOL_SOCKET, MAGMA_NETWORK)
        # before the service is safe to start. The gateway address per network
        # is baked into the binary; there is no policy file to drop in.
        if [ -d /run/systemd/system ]; then
            systemctl --system daemon-reload >/dev/null || true
            # Honour operator's prior enable/disable state on upgrades.
            if [ -n "$2" ] && command -v deb-systemd-invoke >/dev/null 2>&1; then
                deb-systemd-invoke try-restart "${UNIT}" >/dev/null || true
            fi
        fi

        # First-install hint (`$2` empty == fresh install per dpkg(1)).
        if [ -z "$2" ]; then
            cat <<EOF

magma-sidecar installed.

The service is **not** enabled yet. Before starting it:

  1. Edit /etc/magma-sidecar/sidecar.env, at minimum set:
       - MAGMA_MONAD_RPC_URL
       - MAGMA_TXPOOL_SOCKET
       - MAGMA_NETWORK   (mainnet | testnet | localnet)
     The gateway address for each network is baked into the binary; switching
     networks is just changing MAGMA_NETWORK and restarting the service.
  2. systemctl enable --now magma-sidecar
  3. journalctl -u magma-sidecar -f

The service runs as User=monad (matching the validator's IPC socket owner).
EOF
        fi
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        : # nothing to do
        ;;

    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

exit 0
